# conf/.cron
5 * * * * __APP__ /usr/bin/php__PHP_VERSION__ __INSTALL_DIR__/tools/runScheduledTasks.php > /dev/null
